# meta-yttrium-app/recipes-yttrium/yttrium-app/files/yttrium.service
[Unit]
Description=Yttrium demo application

[Service]
ExecStart=@BINDIR@/yttrium.sh

[Install]
WantedBy=multi-user.target
